# supervisord.conf
[supervisord]
nodaemon=true

[program:pixelmap]
command=/usr/bin/python -u /var/www/html/pixelmap_service.py
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true

[program:apache2]
command=/usr/sbin/apache2ctl -e info -DFOREGROUND
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true

[program:parity]
command=parity --geth --datadir /root/ethereum --jsonrpc-hosts all --jsonrpc-interface all --warp
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true

[program:redis]
command=redis-server
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
